# siteeye.service
[Unit]
Description=SiteEye v2 — Wearable AI Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi-molt
WorkingDirectory=/home/pi-molt/molt-device
ExecStart=/home/pi-molt/molt-device/venv/bin/python3 main_v2.py
Restart=on-failure
RestartSec=5
Environment=SITEEYE_PROXY=https://molted.tail4a98c5.ts.net
Environment=TELEGRAM_BOT_TOKEN=8358560979:AAHEPmk-qQg9RsmrIR2cXyNP4i4u_Hqtl2U
Environment=TELEGRAM_CHAT_ID=8217278203

[Install]
WantedBy=multi-user.target
